# portal3-hub.service
[Unit]
Description=Portal3 hub
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /root/portal3-hub/index.js
Restart=on-failure

[Install]
WantedBy=multi-user.target